# lienol/luci-app-koolddns/root/etc/uci-defaults/luci-app-koolddns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@koolddns[-1]
	add ucitrack koolddns
	set ucitrack.@koolddns[-1].init=koolddns
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
