# luci-app-passwall/root/etc/hotplug.d/iface/98-passwall
#!/bin/sh

[[ "$ACTION" == "ifup" && $(uci get "passwall.@global[0].enabled") == "1" ]] && {
    /etc/init.d/passwall restart
    echo "passwall: restart when $INTERFACE ifup" > /dev/kmsg
}
