# lienol/luci-theme-argon-mc/root/etc/uci-defaults/30_luci-theme-argon-mc
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon_MC1=/luci-static/argon_mc1
	set luci.themes.Argon_MC2=/luci-static/argon_mc2
	commit luci
EOF
exit 0
