# lienol/luci-theme-argon-mc/root/etc/uci-defaults/30_luci-theme-argon-mc
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon_MC=/luci-static/argon_mc
	commit luci
EOF
exit 0
