# lienol/luci-theme-netgear-mc/root/etc/uci-defaults/30_luci-theme-netgear-mc
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Netgear_MC=/luci-static/netgear_mc
	commit luci
EOF
exit 0
