# lienol/luci-theme-bootstrap-gaypurple/root/etc/uci-defaults/30_luci-theme-bootstrap-gaypurple
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap_GayPurple=/luci-static/bootstrap_gaypurple
	set luci.main.mediaurlbase=/luci-static/bootstrap_gaypurple
	commit luci
EOF
exit 0
