# lienol/luci-theme-argon-dark-mod/root/etc/uci-defaults/30_luci-theme-argon-dark-mod
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon_Dark=/luci-static/argon_dark
	set luci.themes.Argon_Dark_Purple=/luci-static/argon_dark_purple
	commit luci
EOF
exit 0
