# lienol/luci-theme-bootstrap-purple/root/etc/uci-defaults/30_luci-theme-bootstrap-purple
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap_GayPurple=/luci-static/bootstrap_purple
	set luci.main.mediaurlbase=/luci-static/bootstrap_purple
	commit luci
EOF
exit 0
