# lienol/luci-app-netstatus/root/usr/bin/net_check
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 writen by rapistor(明月永在）
# Must keep author's information if you use this file.


start() {
    status=`ps | grep -c net_check`
    [ "$status" -ge "4" ] && echo "another net_check is running,exit "
    [ "$status" -ge "4" ] && exit

    c_enable=$(uci get netstatus.@netstatus[0].enable)
    check_interval=$(uci get netstatus.@netstatus[0].check_interval)

    if [ $c_enable == "0" ]; then
	exit 0
    fi

    while true
    do
		if ping -w 1 -c 1 www.baidu.com || ping -w 1 -c 1 www.163.com
		then
				echo "Network is online, sleep $check_interval S for next check"
				sleep $check_interval
		else
			echo "Network is offline, network is restarting...."
			/etc/init.d/network restart
			echo "sleep $check_interval S for next check"
			sleep $check_interval
		fi
    done
}